# cron.d/check_in
#自动签到
0 0 * * * bit0r cd ~/crawlers/ && ./offline_account.py
0 9 * * * bit0r cd ~/crawlers/ && ./xtu_covid.py
